# files/etc/uci-defaults/999-sqm-settings
#!/bin/sh
uci -q batch <<'EOI'
set sqm.wan=queue
set sqm.wan.enabled='1'
set sqm.wan.interface='wan'
set sqm.wan.script='nss-zk.qos'
set sqm.wan.qdisc='fq_codel'
commit sqm
EOI

exit 0
